# server/linux/iperf3/iperf3.service
[Unit]
Description=iPerf3 server
After=syslog.target network.target

[Service]
ExecStart=/usr/bin/iperf3 -s

[Install]
WantedBy=multi-user.target
